# curator4-6644102438.conf
[program:curator4-6644102438]
command=python3 /home/curator4/curator4-6644102438/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-6644102438/test.err.log
stdout_logfile=/home/curator4/curator4-6644102438/test.out.log